# examples/http/lw-ppocr-http-onnxruntime.service
[Unit]
Description=lw.PPOCR ONNX Runtime HTTP OCR service
After=network.target

[Service]
Type=simple
DynamicUser=yes
WorkingDirectory=/opt/lw-ppocr
Environment=LD_LIBRARY_PATH=/opt/lw-ppocr:/opt/lw-ppocr/runtimes/linux-x64/onnxruntime
ExecStart=/opt/lw-ppocr/lw-ppocr-http-service --console --config /opt/lw-ppocr/http-service.json
Restart=on-failure
RestartSec=3
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
